# systemd/digiplayer-kiosk.service
[Unit]
Description=DigiPlayer Kiosk Mode (Cage/Wayland)
After=digiplayer.service
Wants=digiplayer.service
Conflicts=getty@tty1.service

[Service]
Type=simple
User=pi
PAMName=login
TTYPath=/dev/tty1
StandardInput=tty
StandardOutput=tty
StandardError=journal
UtmpIdentifier=tty1
UtmpMode=user

# Use cage (Wayland kiosk compositor) with chromium
ExecStart=/usr/bin/cage -s -- /usr/bin/chromium-browser --enable-features=UseOzonePlatform --ozone-platform=wayland --kiosk --noerrdialogs --disable-infobars --no-first-run http://localhost:8080

Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
