# systemd/digiplayer.service
[Unit]
Description=DigiPlayer Digital Signage Client
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 -m digiplayer.main
WorkingDirectory=/opt/digiplayer
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment
Environment=PYTHONUNBUFFERED=1

# Security
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
