# systemd/digiplayer.service
[Unit]
Description=DigiPlayer Main Service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
ExecStart=/usr/bin/python3 -m digiplayer.main --service
WorkingDirectory=/opt/digiplayer
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target
